# terraform/mysql/modules/instance/templates/etc/systemd/system/mysql-replica-dump.timer.tpl
[Unit]
Description=Perform replica sourced, MySQL logical backups according to a schedule

[Timer]
OnCalendar=${schedule}
Unit=mysql-replica-dump.service
Persistent=true

[Install]
WantedBy=timers.target
