# terraform/mysql/modules/instance/templates/etc/systemd/system/mysql-replica-snapshot.timer.tpl
[Unit]
Description=Perform replica sourced, MySQL physical backups according to a schedule

[Timer]
OnCalendar=${schedule}
Unit=mysql-replica-snapshot.service
Persistent=true

[Install]
WantedBy=timers.target
